# php-worker/supervisord.d/hyperf.conf.example
[program:hyperf-worker]
process_name=%(program_name)s_%(process_num)02d
command=php bin/hyperf.php start
autostart=true
autorestart=true
numprocs=8
redirect_stderr=true
